# supervisord.conf
[supervisord]
nodaemon=true
logfile=%(here)s/log/supervisord.log
pidfile=%(here)s/log/supervisord.pid
loglevel=debug

[inet_http_server]
port=127.0.0.1:5700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:config]
command=%(here)s/node_modules/.bin/coffee %(here)s/index.coffee.md
environment=DEBUG="thinkable-ducks:*,useful-wind:*" SUPERVISOR="http://127.0.0.1:5700" MODE=config FSCONF=%(here)s/conf/freeswitch.xml
priority=10

[program:server]
command=%(here)s/node_modules/.bin/coffee %(here)s/index.coffee.md
environment=DEBUG="thinkable-ducks:*,useful-wind:*" SUPERVISOR="http://127.0.0.1:5700" MODE=server
priority=20
autostart=false
autorestart=true

[program:freeswitch]
command=/opt/freeswitch/bin/freeswitch -c -nonat -nonatmap -conf %(here)s/conf -log %(here)s/log -db /dev/shm/freeswitch
autostart=false
autorestart=true
priority=50
